# deploy/rtl-demo-domain-supervisor.conf
[program:redhawk-rtl-demo-domain]
command=/var/redhawk/web/rtl-demo/app/bin/startdomain.sh         ; the program (relative uses PATH, can take args)
user=redhawk
;priority=999                ; the relative start priority (default 999)
autostart=true              ; start at supervisord start (default: true)
autorestart=true            ; retstart at unexpected quit (default: true)
stdout_logfile=/var/log/redhawk-web/rtl-demo-domain.log        ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=1MB   ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10     ; # of stdout logfile backups (default 10)
stderr_logfile=/var/log/redhawk-web/rtl-demo-domain-error.log        ; stderr log path, NONE for none; default AUTO
stderr_logfile_maxbytes=1MB   ; max # logfile bytes b4 rotation (default 50MB)
stderr_logfile_backups=10     ; # of stderr logfile backups (default 10)

